# cron/check_running.crontab
* * * * * /Users/dbuser/bigg_server/bigg_models/bin/check_running.sh 2>> /Users/dbuser/bigg_server/bigg_models.log
